# overpass-slackbuild/rc.local
#!/bin/bash
#
# include the following lines in your /etc/rc.d/rc.local file to start overpass
# dispatcher daemon on boot. Edit "rc.dispatcher" file for your actual path
# to your database directory.

# Startup overpass dispatcher
if [ -x /etc/rc.d/rc.dispatcher ]; then
    /etc/rc.d/rc.dispatcher start
fi
